# labs/guest-session-20161031154351/skel/.config/autostart/screenlocking.desktop
[Desktop Entry]
Name='Disable screen locking in guest session'
Type=Application
Exec=gsettings set org.gnome.desktop.lockdown disable-lock-screen true
